# heroku/supervisor/ss.ini
[program:ss]
command=ss-server -k SS_PASS -p 18081 -m chacha20 -u --plugin v2ray-plugin --plugin-opts "server;host=example.com;path=SS_PATH"
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/ss.log

