# heroku/supervisor/ttyd.ini
[program:ttyd]
command=/usr/sbin/ttyd -p 7681 -i lo sh
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/ttyd.log
