# openshift/supervisor/squid.ini
[program:squid]
command=squid -N -f /etc/squid/squid.conf
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/squid.log
