# openshift/supervisor/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
