# matrix/supervisor/chisel.ini
[program:chisel]
command=/usr/local/bin/chisel server --port 3000 --proxy http://example.com
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/chisel.log
