# matrix/supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon on; master_process on;'
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/nginx.log
