# arukas/supervisor/ss.ini
[program:ss]
command=/usr/local/bin/shadowsocks-server -c /etc/shadowsocks.json
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/ss.log
