# ezgoo/supervisor/ezgoo.ini
[program:ss]
command=/usr/local/bin/ezgoo -dir=/var/www/dist
autostart=true
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/ezgoo.log
